# infra/systemd/kovi-api.service
[Unit]
Description=Kovi API Service
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/kovi-engine
EnvironmentFile=/etc/kovi/kovi.env
ExecStart=/usr/bin/node --import tsx /opt/kovi-engine/apps/api/src/index.ts
Restart=always
RestartSec=5
User=kovi
Group=kovi

[Install]
WantedBy=multi-user.target
